# evoke.service
# Documentation: See evoke.md for usage, configuration, and troubleshooting

[Unit]
Description=Run evoke script on demand

[Service]

Type=exec

ExecStart=/home/danna/endo/evoke.sh

# After script completes, service is marked as inactive (allows repeated runs)
RemainAfterExit=no

# Redirect stdout to systemd journal for easier viewing
StandardOutput=journal

# Redirect stderr to systemd journal for easier viewing
StandardError=journal

